# collectd/rootfs/etc/services.d/collectd/run
#!/usr/bin/with-contenv bashio
# vim: ft=bash
# shellcheck shell=bash
# ==============================================================================
# Start collectd service
# ==============================================================================

bashio::log.info "Service starting"
exec /usr/sbin/collectd -f &

ENABLE_MQTT_PLUGIN=$(bashio::config 'enable_mqtt_plugin')
if [ "${ENABLE_MQTT_PLUGIN}" == "true" ]
then
    bashio::log.info "Creating MQTT sensors"
    /opt/bin/mqtt_sensors.sh
fi  
