# collectd/rootfs/etc/services.d/collectd/run
#!/usr/bin/with-contenv bashio
# vim: ft=bash
# shellcheck shell=bash
# ==============================================================================
# Start collectd service
# ==============================================================================

bashio::log.info "Service starting"
exec /usr/sbin/collectd -f
